# mcp-proxy/rootfs/etc/services.d/mcp-proxy/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Take down the s6 supervision tree when the service exits unexpectedly
# ==============================================================================

declare exit_code=${1}

if [[ "${exit_code}" -ne 0 ]] && [[ "${exit_code}" -ne 256 ]]; then
    bashio::log.error "Service crashed with exit code ${exit_code}, halting..."
    /run/s6/basedir/bin/halt
fi
